# files/anitya.crontab
* * * * * root /src/runcronjobincontainer.sh >> /var/tmp/anitya_cronscript.log 2>&1
